# %HOME%/.config/systemd/user/update-cgit-repos.timer
[Unit]
Description=Run pull.sh periodically

[Timer]
Persistent=false
OnCalendar=daily
AccuracySec=1

[Install]
WantedBy=timers.target
