# systemd/cgitize.timer
[Unit]
Description=Pull cgit repositories periodically

[Timer]
Persistent=false
OnCalendar=daily
AccuracySec=1

[Install]
WantedBy=timers.target
